# config/service/11.smtpd/run
#!/bin/ash
# BoxLinux smtpd service
echo "Mail server"
exec smtpd -d 2>&1
