# sysconfig/service/12.tor/run
#!/bin/ash
# BoxLinux tor anonymity service
echo "Starting tor networking"
exec tor -f /etc/torrc 2>&1
