# sysconfig/service/99.console/run
#!/bin/ash
# BoxLinux getty tty service
echo "$(date) Running getty on tty1"
getty 115200 console linux  
