# sysconfig/service/04.ntpd/run
#!/bin/ash
# BoxLinux network time service
echo "$(date) NTP"
exec ntpd -n -d  2>&1  
