# sysconfig/service/11.smtpd/run
#!/bin/ash
# BoxLinux smtpd service
echo "$(date) Mail server"
exec smtpd -d 2>&1
