# sysconfig/service/03.crond/run
#!/bin/ash
# BoxLinux cron service
echo "Task scheduler"
exec /usr/sbin/crond -f 2>&1  
