# config/service/04.ntpd/run
#!/bin/ash
# BoxLinux network time service
echo "NTP"
exec ntpd -n -d  2>&1  
