# sysconfig/service/02.syslogd/run
#!/bin/ash
# BoxLinux syslog service
echo "System logger"
exec syslogd -n 2>&1
