# config/service/01.klogd/run
#!/bin/ash
# BoxLinux klog service
echo "Kernel logger"
exec klogd -n 2>&1

