# sysconfig/service/01.klogd/run
#!/bin/ash
# BoxLinux klog service
echo "$(date) Kernel logger"
exec klogd -n 2>&1

