# sysconfig/service/12.sshd/run
#!/bin/ash
# BoxLinux sshd service
echo "$(date) SSH server"
exec /usr/sbin/sshd -D 2>&1
