# sysconfig/service/10.httpd/run
#!/bin/ash
# BoxLinux httpd service
echo "$(date) Web server"
exec httpd -u www:www -f -h /srv/www 2>&1
